# pm2-ai.service
[Unit]
Description=PM2 Process Manager for AI Runlevel
After=network.target

[Service]
Type=forking
User=loop
LimitNOFILE=infinity
LimitNPROC=infinity
LimitCORE=infinity
Environment=PATH=/usr/bin:/usr/local/bin:/home/loop/.npm-global/bin
Environment=PM2_HOME=/home/loop/.pm2
ExecStart=/usr/local/bin/pm2 start /home/loop/ai-runlevel/ecosystem.config.cjs
ExecReload=/usr/local/bin/pm2 reload all
ExecStop=/usr/local/bin/pm2 kill

[Install]
WantedBy=multi-user.target